# pnetlab-schema/xcontrol/postinst
#!/bin/bash
systemctl enable mysql
service mysql start
exit 0